# rust-mcp/debian/postinst
#!/bin/bash
set -e

# Post-installation script for rust-mcp
# Creates user configuration directory and default env file

# Get the actual user (not root when using sudo)
if [ -n "$SUDO_USER" ]; then
    ACTUAL_USER="$SUDO_USER"
    ACTUAL_HOME=$(getent passwd "$SUDO_USER" | cut -d: -f6)
else
    ACTUAL_USER="$USER"
    ACTUAL_HOME="$HOME"
fi

# Skip if no valid user found (e.g., during automated builds)
if [ -z "$ACTUAL_USER" ] || [ "$ACTUAL_USER" = "root" ]; then
    echo "Note: Skipping user config setup (no regular user detected)"
    echo "Run 'rust-mcp-service --help' as your user to create config"
    exit 0
fi

CONFIG_DIR="$ACTUAL_HOME/.config/rust-mcp"

# Create config directory if it doesn't exist
if [ ! -d "$CONFIG_DIR" ]; then
    mkdir -p "$CONFIG_DIR"
    chown "$ACTUAL_USER:$ACTUAL_USER" "$CONFIG_DIR"
    echo "Created config directory: $CONFIG_DIR"
fi

# Create instances.json for multi-instance configuration
if [ ! -f "$CONFIG_DIR/instances.json" ]; then
    cat > "$CONFIG_DIR/instances.json" << 'EOF'
{
  "production": {
    "url": "http://localhost:8069",
    "db": "production",
    "apiKey": "YOUR_ODOO_19_API_KEY"
  },
  "staging": {
    "url": "http://localhost:8069",
    "db": "staging",
    "apiKey": "YOUR_STAGING_API_KEY"
  },
  "development": {
    "url": "http://localhost:8069",
    "db": "development",
    "version": "18",
    "username": "admin",
    "password": "admin"
  }
}
EOF
    chown "$ACTUAL_USER:$ACTUAL_USER" "$CONFIG_DIR/instances.json"
    chmod 600 "$CONFIG_DIR/instances.json"
    echo "Created default instances.json: $CONFIG_DIR/instances.json"
fi

# Create default env file if it doesn't exist
if [ ! -f "$CONFIG_DIR/env" ]; then
    cat > "$CONFIG_DIR/env" << EOF
# Odoo Rust MCP Server Configuration
# Multi-instance configuration (default)

# =============================================================================
# Multi-Instance Configuration (Default - Recommended)
# =============================================================================
# Uses instances.json for multiple Odoo instances
ODOO_INSTANCES_JSON=$CONFIG_DIR/instances.json

# =============================================================================
# Single Instance Configuration (Alternative - uncomment if not using multi-instance)
# =============================================================================
# # Odoo 19+ (API Key authentication)
# ODOO_URL=http://localhost:8069
# ODOO_DB=mydb
# ODOO_API_KEY=YOUR_API_KEY
#
# # Odoo < 19 (Username/Password authentication)
# # ODOO_VERSION=18
# # ODOO_USERNAME=admin
# # ODOO_PASSWORD=admin

# =============================================================================
# Config UI Authentication
# =============================================================================
# Login credentials for the config web UI
# IMPORTANT: Change these default credentials!
CONFIG_UI_USERNAME=admin
CONFIG_UI_PASSWORD=changeme

# =============================================================================
# MCP HTTP Transport Authentication
# =============================================================================
# Enable/disable HTTP authentication (default: false)
MCP_AUTH_ENABLED=false
# Generate a secure token: openssl rand -hex 32
# MCP_AUTH_TOKEN=your-secure-random-token-here

# =============================================================================
# MCP Config Paths
# =============================================================================
# Path to MCP configuration files (tools, prompts, server settings)
MCP_TOOLS_JSON=$CONFIG_DIR/tools.json
MCP_PROMPTS_JSON=$CONFIG_DIR/prompts.json
MCP_SERVER_JSON=$CONFIG_DIR/server.json
EOF
    chown "$ACTUAL_USER:$ACTUAL_USER" "$CONFIG_DIR/env"
    chmod 600 "$CONFIG_DIR/env"
    echo "Created default env file: $CONFIG_DIR/env"
else
    # Migration: Add MCP config paths to existing env file if missing
    if ! grep -q "MCP_TOOLS_JSON" "$CONFIG_DIR/env"; then
        echo "Migration: Adding MCP config paths to env file..."
        cat >> "$CONFIG_DIR/env" << EOF

# =============================================================================
# MCP Config Paths (added in v0.3.27)
# =============================================================================
# Path to MCP configuration files (tools, prompts, server settings)
MCP_TOOLS_JSON=$CONFIG_DIR/tools.json
MCP_PROMPTS_JSON=$CONFIG_DIR/prompts.json
MCP_SERVER_JSON=$CONFIG_DIR/server.json
EOF
        echo "Migration complete"
    fi
fi

echo ""
echo "=============================================="
echo "rust-mcp installed successfully!"
echo "=============================================="
echo ""
echo "Configuration:"
echo "  Instances: $CONFIG_DIR/instances.json"
echo "  Env file:  $CONFIG_DIR/env"
echo ""
echo "Quick start:"
echo "  1. Edit instances: nano $CONFIG_DIR/instances.json"
echo "  2. Start service:  sudo systemctl start rust-mcp"
echo "  3. Enable on boot: sudo systemctl enable rust-mcp"
echo ""
echo "Service endpoint: http://127.0.0.1:8787/mcp"
echo ""

#DEBHELPER#
