# rust-mcp/debian/rust-mcp.service
[Unit]
Description=Odoo MCP Server - Model Context Protocol server for Odoo integration
Documentation=https://github.com/rachmataditiya/odoo-rust-mcp
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/rust-mcp-service --transport http --listen 127.0.0.1:8787
Restart=on-failure
RestartSec=5

# Run as user service when possible, fallback to root
# For user-level service, use: systemctl --user start rust-mcp
DynamicUser=no

# Security hardening (relaxed to allow config file access)
NoNewPrivileges=true
ProtectSystem=full
PrivateTmp=true

# Allow read-write access to user config directories
# The service wrapper creates config in ~/.config/rust-mcp/
ReadWritePaths=/home
ReadWritePaths=/root/.config

# Working directory for static files
WorkingDirectory=/usr/share/rust-mcp

# Environment
Environment=RUST_LOG=info

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=rust-mcp

[Install]
WantedBy=multi-user.target
